# .config/systemd/user/openclaw-gateway.service
[Unit]
Environment=SONARR_URL=http://localhost:8989
Environment=SONARR_API_KEY=dfa7da413b48480097f6a399080be393
Environment=SONARR_URL=http://localhost:8989
Environment=OPENAI_API_KEY=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
Environment=SONARR_URL=http://localhost:8989
Environment=SONARR_API_KEY=dfa7da413b48480097f6a399080be393
Environment=SONARR_URL=http://localhost:8989
Description=OpenClaw Gateway (v2026.2.21-2)
Environment=SONARR_URL=http://localhost:8989
Environment=SONARR_API_KEY=dfa7da413b48480097f6a399080be393
Environment=SONARR_URL=http://localhost:8989
Environment=OPENAI_API_KEY=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
Environment=SONARR_URL=http://localhost:8989
Environment=SONARR_API_KEY=dfa7da413b48480097f6a399080be393
Environment=SONARR_URL=http://localhost:8989
After=network-online.target
Environment=SONARR_URL=http://localhost:8989
Environment=SONARR_API_KEY=dfa7da413b48480097f6a399080be393
Environment=SONARR_URL=http://localhost:8989
Environment=OPENAI_API_KEY=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
Environment=SONARR_URL=http://localhost:8989
Environment=SONARR_API_KEY=dfa7da413b48480097f6a399080be393
Environment=SONARR_URL=http://localhost:8989
Wants=network-online.target
Environment=SONARR_URL=http://localhost:8989
Environment=SONARR_API_KEY=dfa7da413b48480097f6a399080be393
Environment=SONARR_URL=http://localhost:8989
Environment=OPENAI_API_KEY=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
Environment=SONARR_URL=http://localhost:8989
Environment=SONARR_API_KEY=dfa7da413b48480097f6a399080be393
Environment=SONARR_URL=http://localhost:8989

[Service]
Environment=SONARR_URL=http://localhost:8989
Environment=SONARR_API_KEY=dfa7da413b48480097f6a399080be393
Environment=SONARR_URL=http://localhost:8989
Environment=OPENAI_API_KEY=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
Environment=SONARR_URL=http://localhost:8989
Environment=SONARR_API_KEY=dfa7da413b48480097f6a399080be393
Environment=SONARR_URL=http://localhost:8989
ExecStart=/home/jules/.nvm/versions/node/v22.22.0/bin/node /home/jules/.nvm/versions/node/v22.22.0/lib/node_modules/openclaw/dist/index.js gateway --port 18789
Environment=SONARR_URL=http://localhost:8989
Environment=SONARR_API_KEY=dfa7da413b48480097f6a399080be393
Environment=SONARR_URL=http://localhost:8989
Environment=OPENAI_API_KEY=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
Environment=SONARR_URL=http://localhost:8989
Environment=SONARR_API_KEY=dfa7da413b48480097f6a399080be393
Environment=SONARR_URL=http://localhost:8989
Restart=always
Environment=SONARR_URL=http://localhost:8989
Environment=SONARR_API_KEY=dfa7da413b48480097f6a399080be393
Environment=SONARR_URL=http://localhost:8989
Environment=OPENAI_API_KEY=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
Environment=SONARR_URL=http://localhost:8989
Environment=SONARR_API_KEY=dfa7da413b48480097f6a399080be393
Environment=SONARR_URL=http://localhost:8989
RestartSec=5
Environment=SONARR_URL=http://localhost:8989
Environment=SONARR_API_KEY=dfa7da413b48480097f6a399080be393
Environment=SONARR_URL=http://localhost:8989
Environment=OPENAI_API_KEY=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
Environment=SONARR_URL=http://localhost:8989
Environment=SONARR_API_KEY=dfa7da413b48480097f6a399080be393
Environment=SONARR_URL=http://localhost:8989
KillMode=process
Environment=SONARR_URL=http://localhost:8989
Environment=SONARR_API_KEY=dfa7da413b48480097f6a399080be393
Environment=SONARR_URL=http://localhost:8989
Environment=OPENAI_API_KEY=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
Environment=SONARR_URL=http://localhost:8989
Environment=SONARR_API_KEY=dfa7da413b48480097f6a399080be393
Environment=SONARR_URL=http://localhost:8989
Environment=HOME=/home/jules
Environment=SONARR_URL=http://localhost:8989
Environment=SONARR_API_KEY=dfa7da413b48480097f6a399080be393
Environment=SONARR_URL=http://localhost:8989
Environment=OPENAI_API_KEY=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
Environment=SONARR_URL=http://localhost:8989
Environment=SONARR_API_KEY=dfa7da413b48480097f6a399080be393
Environment=SONARR_URL=http://localhost:8989
Environment=TMPDIR=/tmp
Environment=SONARR_URL=http://localhost:8989
Environment=SONARR_API_KEY=dfa7da413b48480097f6a399080be393
Environment=SONARR_URL=http://localhost:8989
Environment=OPENAI_API_KEY=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
Environment=SONARR_URL=http://localhost:8989
Environment=SONARR_API_KEY=dfa7da413b48480097f6a399080be393
Environment=SONARR_URL=http://localhost:8989
Environment=PATH=/home/jules/.local/bin:/home/jules/.npm-global/bin:/home/jules/bin:/home/jules/.volta/bin:/home/jules/.asdf/shims:/home/jules/.bun/bin:/home/jules/.nvm/current/bin:/home/jules/.fnm/current/bin:/home/jules/.local/share/pnpm:/usr/local/bin:/usr/bin:/bin
Environment=SONARR_URL=http://localhost:8989
Environment=SONARR_API_KEY=dfa7da413b48480097f6a399080be393
Environment=SONARR_URL=http://localhost:8989
Environment=OPENAI_API_KEY=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
Environment=SONARR_URL=http://localhost:8989
Environment=SONARR_API_KEY=dfa7da413b48480097f6a399080be393
Environment=SONARR_URL=http://localhost:8989
Environment=OPENCLAW_GATEWAY_PORT=18789
Environment=SONARR_URL=http://localhost:8989
Environment=SONARR_API_KEY=dfa7da413b48480097f6a399080be393
Environment=SONARR_URL=http://localhost:8989
Environment=OPENAI_API_KEY=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
Environment=SONARR_URL=http://localhost:8989
Environment=SONARR_API_KEY=dfa7da413b48480097f6a399080be393
Environment=SONARR_URL=http://localhost:8989
Environment=OPENCLAW_SYSTEMD_UNIT=openclaw-gateway.service
Environment=SONARR_URL=http://localhost:8989
Environment=SONARR_API_KEY=dfa7da413b48480097f6a399080be393
Environment=SONARR_URL=http://localhost:8989
Environment=OPENAI_API_KEY=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
Environment=SONARR_URL=http://localhost:8989
Environment=SONARR_API_KEY=dfa7da413b48480097f6a399080be393
Environment=SONARR_URL=http://localhost:8989
Environment=OPENCLAW_SERVICE_MARKER=openclaw
Environment=SONARR_URL=http://localhost:8989
Environment=SONARR_API_KEY=dfa7da413b48480097f6a399080be393
Environment=SONARR_URL=http://localhost:8989
Environment=OPENAI_API_KEY=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
Environment=SONARR_URL=http://localhost:8989
Environment=SONARR_API_KEY=dfa7da413b48480097f6a399080be393
Environment=SONARR_URL=http://localhost:8989
Environment=OPENCLAW_SERVICE_KIND=gateway
Environment=SONARR_URL=http://localhost:8989
Environment=SONARR_API_KEY=dfa7da413b48480097f6a399080be393
Environment=SONARR_URL=http://localhost:8989
Environment=OPENAI_API_KEY=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
Environment=SONARR_URL=http://localhost:8989
Environment=SONARR_API_KEY=dfa7da413b48480097f6a399080be393
Environment=SONARR_URL=http://localhost:8989
Environment=OPENCLAW_SERVICE_VERSION=2026.2.21-2
Environment=SONARR_URL=http://localhost:8989
Environment=SONARR_API_KEY=dfa7da413b48480097f6a399080be393
Environment=SONARR_URL=http://localhost:8989
Environment=OPENAI_API_KEY=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
Environment=SONARR_URL=http://localhost:8989
Environment=SONARR_API_KEY=dfa7da413b48480097f6a399080be393
Environment=SONARR_URL=http://localhost:8989

[Install]
WantedBy=default.target
Environment=SONARR_URL=http://localhost:8989
Environment=SONARR_API_KEY=dfa7da413b48480097f6a399080be393
Environment=SONARR_URL=http://localhost:8989
Environment=OPENAI_API_KEY=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
Environment=SONARR_URL=http://localhost:8989
Environment=SONARR_API_KEY=dfa7da413b48480097f6a399080be393
Environment=SONARR_URL=http://localhost:8989
