# .config/systemd/user/pi-bot.service
[Unit]
Description=Raspberry Pi 5 Media Bot (Deterministic)
After=network.target

[Service]
ExecStart=/usr/bin/python3 /home/jules/scripts/pi_bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
